# calendar.cron
1 1 * * * /home/mike/2576-lwb/update_calendar.sh > /home/mike/2576-lwb/output/calendar_update.log 2>&1

